# ss.conf
[program:ss]
command=sudo ssserver -c config.json
directory=/etc/shadowsocks
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
